# wireless/staassochistory/new_wireless_tools/wireless_tools_v0.2_20170228/spf_openwrt_rulepkg/wireless-tools/files/iwevent.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org


START=80

BIN=iwevent

start() {
	echo -e "\033[32mSTART THE SERVICE\033[0m"
	$BIN  >>/dev/null 2>&1 &
	
}


stop(){
	killall iwevent

}
