# wireless/skysoft_ieee1905/files/ieee1905.init
#!/bin/sh /etc/rc.common


start(){
 echo "start" >/dev/console

}

stop(){
echo "stop" >/dev/console


}
restart(){
	stop
	start

}

