# wireless/proc-reffer/skysoft_detection/files/mail.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010 OpenWrt.org

#START=90

#SERVICE_DAEMONIZE=1
APP_BIN="/usr/sbin/mymail"

start() {
        service_start $APP_BIN
}

stop() {
        service_stop $APP_BIN
        rm /var/my.txt
}
