# services/test/netinfo_notify/files/netinfo_notify.init
#!/bin/sh /etc/rc.common


BIN="/sbin/netinfo_notify"

start()
{
	daemon=`basename $BIN`
	
	$BIN
	 
	prog=`ps | grep $daemon | grep -v "grep"`
	
	if [ -n "$prog" ];then
		echo -e "stat the $daemon [\033[32m SUCCEED \033[0m]"
	fi
}

stop(){
	killall -9 netinfo_notify
	
}
